# .config/profile.d/aliases
alias b='bbedit .'
alias be='bundle exec'
alias bi='bundle install'
alias bu='bundle update --source'
alias cp='cp -i'
alias d='cd spec/dummy'
alias findip='ifconfig | grep -A 6 "en" | grep "inet "'
alias ga='git add .'
alias gb='git branch'
alias gc='git commit'
alias gcb='git checkout -b'
alias gcm='git checkout master'
alias gd='git diff'
alias gl='git log'
alias gp='git push -u origin $(git branch | grep "^*" | awk "{ print \$2 }")'
alias gs='git status'
alias gu='git fetch -p ; git rebase ; git branch --merged | grep -v master | grep -v "^\\*" | xargs git branch -d'
alias ls='/bin/ls -FG'
alias l='ls'
alias la='ls -a'
alias less='less -R'
alias ll='ls -l'
alias lla='ls -a -l'
alias logd='tail -F log/development.log'
alias logp='tail -F log/production.log'
alias logs='tail -F log/*.log'
alias logt='tail -F log/test.log'
alias mv='mv -i'
alias p='cd -'
alias psgrep="ps aux | grep -v grep | grep"
alias r='touch tmp/restart.txt'
alias re="source $HOME/.bash_profile"
alias rm='rm -i'
alias s='cd .. && pwd'
