# .config/profile.d/aliases
alias c='cd $_'
alias cp='cp -i'
alias ls='/bin/ls -FG'
alias l='ls'
alias la='ls -a'
alias less='less -R'
alias ll='ls -l'
alias lla='ls -a -l'
alias mv='mv -i'
alias p='cd -'
alias rm='rm -i'
alias s='cd .. && pwd'
alias v='vim $_'

export PATH="$HOME/bin:$PATH"
