# .config/profile.d/bash
alias cp='cp -i'
alias findip='ifconfig | grep -A 6 "en" | grep "inet "'
alias ls='/bin/ls -FG'
alias l='ls'
alias la='ls -a'
alias less='less -R'
alias ll='ls -l'
alias lla='ls -a -l'
alias mv='mv -i'
alias p='cd -'
alias psgrep="ps aux | grep -v grep | grep"
alias rm='rm -i'
alias s='cd .. && pwd'
