# .config/profile.d/completions.ruby
complete -C $HOME/bin/completions -o default rake
complete -C $HOME/bin/completions -o default rails
complete -C $HOME/bin/completions -o default cap
