# .config/profile.d/rails
alias d='test -d test && cd test/dummy || cd spec/dummy'
alias logd='tail -F log/development.log'
alias logp='tail -F log/production.log'
alias logs='tail -F log/*.log'
alias logt='tail -F log/test.log'
alias vmig='vim db/migrate/$(ls -1 db/migrate | tail -1)'
alias vimg='vmig'

routes(){
  rake routes | grep "$1"
}
