# .config/profile.d/ruby
gemcd(){
  if [ -n "$1" ]; then
    cd $(bundle info $1 | perl -lne 'print $1 if /Path:\s*(.*)/')
  else
    cd $(dirname $(bundle list --paths | tail -n 1))
  fi

  pwd
}
