# .config/profile.d/yarn
alias yi='yarn install --check-files'
alias yin='yarn init'
